# netsonar.service
[Unit]
Description=NetSonar - Prometheus network monitoring exporter
Documentation=https://github.com/luk-kop/netsonar
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=netsonar
Group=netsonar

# Optional overrides (e.g. NETSONAR_OPTS="--listen-addr :9999").
# The leading dash means "don't fail if the file doesn't exist".
EnvironmentFile=-/etc/default/netsonar

ExecStart=/usr/local/bin/netsonar --config /etc/netsonar/config.yaml $NETSONAR_OPTS
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s

SyslogIdentifier=netsonar
StandardOutput=journal
StandardError=journal

# Grant CAP_NET_RAW for MTU probes (raw ICMP sockets). Regular ICMP probes use
# unprivileged udp4 sockets and do not require this capability. Operators using
# only TCP/HTTP/DNS/TLS/Proxy probes can safely remove AmbientCapabilities and
# CapabilityBoundingSet.
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ReadOnlyPaths=/etc/netsonar
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
MemoryDenyWriteExecute=true
LockPersonality=true
SystemCallArchitectures=native

# Resource limits
MemoryMax=256M
CPUQuota=25%

[Install]
WantedBy=multi-user.target
